# .crontab
* * * * * php /Applications/MAMP/htdocs/berserk_3/test.php -new_scan -s -move
